# install/etc/xdg/autostart/user-mount-check.desktop
[Desktop Entry]
Name=User mount checker
Comment=Script to check the status of user network directories
Exec=/usr/bin/user-mount-check
Icon=user-mount-check
Terminal=false
Type=Application
Categories=Settings;System;
X-KDE-autostart-phase=2
X-KDE-autostart-after=krunner
